# contrib/nghttpx.service.in
[Unit]
Description=HTTP/2 proxy
After=network.target

[Service]
Type=forking
ExecStart=@bindir@/nghttpx --conf=/etc/nghttpx/nghttpx.conf --pid-file=/run/nghttpx.pid --daemon

[Install]
WantedBy=multi-user.target
